# setup/splash/splash-tap-watcher.service
[Unit]
Description=Reveal boot logs when the splash screen is tapped
DefaultDependencies=no
After=local-fs.target
ConditionKernelCommandLine=splash

[Service]
Type=simple
ExecStart=/usr/local/sbin/splash-tap-watcher
TimeoutStartSec=0

[Install]
WantedBy=sysinit.target
