# kilodash.service
[Unit]
Description=Scottina touchscreen control panel
After=network-online.target NetworkManager.service
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/kilodash
ExecStart=/usr/bin/python3 /opt/kilodash/run.py
Restart=always
RestartSec=2
# The app owns /dev/fb0 directly; keep systemd from wiring a tty to it.
StandardInput=null
StandardOutput=journal
StandardError=journal
# The ILI9486 panel only scans out over SPI when its DRM pipe is enabled.
# fbcon enables it — but only if the console is mapped onto fb0. Force that
# here (survives even if the kernel cmdline maps the console elsewhere), then
# blank the cursor. The short sleep lets the modeset settle before we draw.
ExecStartPre=-/bin/sh -c 'con2fbmap 1 0; con2fbmap 2 0; echo 0 > /sys/class/graphics/fbcon/cursor_blink 2>/dev/null; echo 1 > /proc/sys/kernel/printk; sleep 0.4'

[Install]
WantedBy=multi-user.target
